# pbp-core/systemd/pbp-integrity.timer
[Unit]
Description=Daily PBP Hull Integrity Scan

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
