# pbp-core/systemd/pbp-watch.timer
[Unit]
Description=Periodic PBP Security Watch

[Timer]
OnUnitActiveSec=15min
OnBootSec=5min

[Install]
WantedBy=timers.target
